# service_files/xyz.openbmc_project.cpuutilization.service
[Unit]
Description=CPU Utilization Monitoring Service
After=phosphor-dbus-interfaces.service

[Service]
ExecStart=/usr/bin/cpuutilization
Restart=always
Type=dbus
BusName=xyz.openbmc_project.CpuUtilization

[Install]
WantedBy=multi-user.target
